# {{cookiecutter.project_slug}}/conf/production/supervisor.conf
[program:{{ cookiecutter.project_user }}_gunicorn]
command=/home/{{ cookiecutter.project_user }}/{{ cookiecutter.project_slug }}/utility/production/gunicorn-start.sh
user={{ cookiecutter.project_user }}
umask=0022
autostart=true
autorestart=true
startsecs=3

stdout_logfile=/home/{{ cookiecutter.project_user }}/{{ cookiecutter.project_slug }}/logs/gunicorn.log
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
; process_name=%(program_name)s_%(process_num)02d

stderr_logfile=/home/{{ cookiecutter.project_user }}/{{ cookiecutter.project_slug }}/logs/%(program_name)s_stderr.log
stdout_logfile=/home/{{ cookiecutter.project_user }}/{{ cookiecutter.project_slug }}/logs/%(program_name)s_stdout.log
redirect_stderr=true
